# roles/webgate/files/shibbolethConfigure.service.d/override.conf
[Service]
ExecStartPre=/usr/bin/consul-template -once '-template=/etc/consul-template.d/shibboleth2.xml.tmpl:/etc/shibboleth/shibboleth2.xml'
