# roles/webgate/files/consul.service.d/override.conf
# [Unit]
# Requires=rsyslog.service network-online.target
# After=rsyslog.service network-online.target

[Service]
ExecStartPre=/usr/bin/setup-peers.py
ExecStartPost=/usr/bin/wait-raft.py