# roles/webgate/files/shibd.service.d/override.conf
[Unit]
After=shibbolethConfigure.service consul2file.service
