# roles/webgate/files/consul2file.service.d/override.conf
[Unit]
After=consul2fileConfigure.service

[Service]
ExecStop=
